# exbtn.init
#!/sbin/openrc-run

pidfile='/var/run/exbtn.pid'
command='/usr/local/sbin/exbtn'
start_stop_daemon_args='--background --make-pidfile'
